# librespot-java/librespot-java-restart.service
[Unit]
Description=Restart librespot-java service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart librespot-java.service
